# templates/gmediarender.service.j2
# generated by ansible for host {{ inventory_hostname }}
[Unit]
Description={{ gmrrs_description }}
After=network.target

[Service]
Type=simple
User={{ gmrrs_user }}
ExecStart={{ gmrrs_run }}
Restart=on-abort

[Install]
WantedBy=multi-user.target
